# systemd/hdmi-disable.service
[Unit]
Description=Disable HDMI
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tvservice -o
